# systemd/i3manager.service
[Unit]
Description=i3manager

[Service]
ExecStart=%h/i3-manager/run.sh
